# lib/ansible_collections/local/uws/roles/dunst/files/temperature.service
[Unit]
Description=Notify on critical hardware temperature

[Service]
Type=oneshot
ExecStart=/opt/dunst/bin/temperature.sh
